# cron/reanalyzeRepos.cron
*/10 * * * * NODE_ENV=production /usr/local/bin/node /usr/src/app/scripts/reanalysis.job.js >> /var/log/cron.log 2>&1
