# cron/refreshDefaultRepo.cron
0 8 */2 * * NODE_ENV=production /usr/local/bin/node /usr/src/app/scripts/defaultRepo/createDefaultRepo.js >> /var/log/cron.log 2>&1
